# deploy/tcast-relay.service
# systemd unit for the tcast relay.
#
# Install:
#   sudo cp target/release/tcast-relay /usr/local/bin/tcast-relay
#   sudo cp deploy/tcast-relay.service /etc/systemd/system/
#   echo 'TCAST_AUTH_KEY=change-me' | sudo tee /etc/tcast-relay.env
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now tcast-relay
#
# The relay binds to localhost only; TLS is terminated by Caddy in front of it.

[Unit]
Description=tcast relay
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Optional: a shared host secret. Remove the line to run an open relay.
EnvironmentFile=-/etc/tcast-relay.env
ExecStart=/usr/local/bin/tcast-relay --bind 127.0.0.1:4455
Restart=on-failure
RestartSec=2
# Hardening
DynamicUser=yes
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
